# deploy/pa-agent.service
[Unit]
Description=Personal AI Agent Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/pa-agent
Environment=PYTHONPATH=/root/pa-agent
Environment=OMP_NUM_THREADS=1
Environment=TOKENIZERS_PARALLELISM=false
Environment=TRANSFORMERS_OFFLINE=1
ExecStart=/root/pa-agent/.venv/bin/uvicorn server.main:app \
    --host 127.0.0.1 \
    --port 8000 \
    --workers 1
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
